# config/includes.chroot_after_packages/etc/profile.d/user-bin.sh
# Add user-local bin directories to PATH
PATH="$HOME/bin:$HOME/.local/bin:$PATH"
export PATH

