# config/includes.chroot_after_packages/etc/systemd/system/ensure-sudo.service
[Unit]
Description=Ensure all normal users are in sudo group
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ensure-sudo.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
