# config/includes.chroot_after_packages/etc/skel/.config/autostart/gxkb-tray-icon.desktop
[Desktop Entry]
Name=gxkb launcher
Type=Application
Exec=gxkb &
Terminal=false
Icon=web-browser
Comment=gxkb tray launcher
X-GNOME-Autostart-enabled=true
Hidden=false
